# roles/ansible-corpsys_users/files/profile.d/editor.sh
## This file is controlled with puppet.
## Any manual changes will be overwritten.

# we prefer vim as default editor
export EDITOR="/bin/vi"
