# roles/ansible-corpsys_users/files/profile.d/path.sh
## This file is controlled with puppet.
## Any manual changes will be overwritten.

export PATH="$PATH:/sbin:/usr/sbin:/usr/local/sbin"
